# deploy/polybot-watchdog.service
[Unit]
Description=Polybot Wallet Watchdog
After=network.target

[Service]
Type=simple
User=clem_garnero753
WorkingDirectory=/home/clem_garnero753/polymarket-bot
EnvironmentFile=/home/clem_garnero753/polymarket-bot/.env
Environment=PYTHONPATH=/home/clem_garnero753/polymarket-bot/src
ExecStart=/home/clem_garnero753/polymarket-bot/.venv/bin/python scripts/run_wallet_watchdog.py
Restart=always
RestartSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
